# systemd/roomba_api.service
[Unit]
Description=Roomba lirc API
After=multi-user.target

[Service]
User=roomba_api
Group=roomba_api
WorkingDirectory=/home/roomba_api/roomba_ir_api/
ExecStart=/home/roomba_api/venv_roomba_api/bin/gunicorn --bind 0.0.0.0:8124 roomba_ir_api.app
Type=idle

[Install]
WantedBy=multi-user.target
